# com.alylmz.notdefteri.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.alylmz.notdefteri.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/alylmztr/Documents/GitHub/alylmz-kisisel-not-defterim/backup_to_icloud.sh</string>
    </array>

    <!-- Her 30 dakikada bir çalış -->
    <key>StartInterval</key>
    <integer>1800</integer>

    <!-- Mac açıldığında çalış -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/notdefteri-backup.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/notdefteri-backup-error.log</string>
</dict>
</plist>
